# applications/luci-app-timecontrol/luci-app-timecontrol/root/etc/uci-defaults/luci-app-timecontrol
#!/bin/sh

/etc/init.rd/rpcd restart

exit 0
